# root/etc/cron.d/ptrack_purge
#
# ptrack_purge: clean ptrack *.dump files older than seven days
#

15 0 * * 0	root	find /var/spool/ptrack -name "*.dump" -type f -mtime +7 -exec rm -f {} \;
